# gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn "flaskr:create_app()" -b 127.0.0.1:5000
directory=/deploy
